# tests/tools/docker/centos7-nodejs-echarts-draw/supervisord-nodejs.ini
[program:node]
command=/usr/bin/node /data/www/echarts-draw/server.js